# photo.timer
[Unit]
Description=Take a photo every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target


