# Ansible/roles/raspicam/templates/raspivid.socket.j2
[Unit]
Description=Raspivid Socket
After=wg-quick@server.service
Requires=wg-quick@server.service

[Socket]
ListenStream={{ wg_ip }}:3333
Accept=yes

[Install]
WantedBy=multi-user.target

